# release/groundseg.service
[Unit]
Description=NativePlanet GroundSeg Controller
After=multi-user.target
[Service]
Type=forking
Restart=always
User=root
WorkingDirectory=/opt/nativeplanet/groundseg/
ExecStartPre=/opt/nativeplanet/groundseg/groundseg.sh
ExecStart=docker-compose up
ExecStop=docker-compose down
[Install]
WantedBy=multi-user.target
