# release/groundseg.service
[Unit]
Description=NativePlanet GroundSeg Controller
After=multi-user.target
[Service]
Type=simple
Restart=always
User=root
WorkingDirectory=/opt/nativeplanet/groundseg/
ExecStart=export HOST_HOSTNAME=$(hostname) && docker-compose up -d
ExecStop=docker-compose down
[Install]
WantedBy=multi-user.target
